# hyprland.service
[Unit]
Type=simple
[Service]
ExecStart=/bin/hyprland
[Install]
WantedBy=default.target
