# Linux/hyprland.service
[Unit]
Type=simple
[Service]
ExecStart=/bin/hyprland
[Install]
WantedBy=graphical.target
